# server-deploy/cookbooks/base/files/default/reload-haproxy.timer
[Unit]
Description=HAProxy Reload Timer

[Timer]
OnBootSec=10sec
OnUnitActiveSec=5min
