# mesos-setup/cookbooks/base/files/default/reload-haproxy.timer
[Unit]
Description=HAProxy Reload Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
